# etc/profile.d/no-vsync.sh
__GL_SYNC_TO_VBLANK=0
vblank_mode=0

export __GL_SYNC_TO_VBLANK vblank_mode
