# etc/profile.d/no-vsync.sh
export __GL_SYNC_TO_VBLANK=0
export vblank_mode=0
