# etc/profile.d/no-wine-links.sh
export WINEDLLOVERRIDES=winemenubuilder.exe,mscoree,mshtml=
