# etc/profile.d/no-wine-links.sh
WINEDLLOVERRIDES=winemenubuilder.exe,mscoree,mshtml=

export WINEDLLOVERRIDES
