# etc/profile.d/vsync.sh
export __GL_SYNC_TO_VBLANK=0
export CLUTTER_VBLANK=none
export vblank_mode=0
